# src/rootfs/etc/init.d/iperf3
#!/bin/sh /etc/rc.common

START=65
STOP=85
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/iperf3 -s
	procd_set_param respawn
	procd_close_instance
}