# src/rootfs/etc/uci-defaults/80-custom-network
#!/bin/sh

uci -q batch << EOF
set network.lan.ipaddr='192.168.200.1'
set network.@device[0].bridge_empty='1'
commit network
EOF
