# src/rootfs/etc/uci-defaults/90-custom-system
#!/bin/sh

uci -q batch << EOF
set system.@system[0].hostname='openwrt-iscsi-target'
set system.@system[0].ttylogin='1'
commit system
EOF
