# src/rootfs/etc/uci-defaults/90-custom-password
#!/bin/sh
set -e

# See https://openwrt.org/docs/guide-user/security/dropbear.public-key.auth to use key-based auth

PASSWORD="pass1234"
echo -e "$PASSWORD\n$PASSWORD" | passwd
